# deploy/oz-relay-server.service
[Unit]
Description=OZ Relay Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=oz-relay
Group=oz-relay
WorkingDirectory=/opt/oz-relay

ExecStart=/opt/oz-relay/bin/oz-relay-server
Restart=always
RestartSec=5

EnvironmentFile=/opt/oz-relay/config/config.env
Environment=HOME=/opt/oz-relay

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/oz-relay /opt/arcflow /opt/promotions
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
MemoryMax=8G

# Logging
StandardOutput=append:/opt/oz-relay/logs/server.log
StandardError=append:/opt/oz-relay/logs/server.log

[Install]
WantedBy=multi-user.target
